# taobao/trunk/bin/buyu.cron
#! /bin/sh
cd $1
$2 -c /usr/local/php-cgi/lib/php.bin.ini -f buyu.php
